# supervisor.conf
[program:sentry]
command=docker-compose up
directory=/opt/sentry/setup
redirect_stderr=true
autostart=true
autorestart=true
priority=10
